# color-the-map.service
[Unit]
Description=Color The Map - GPS Track Visualizer
After=network.target

[Service]
Type=simple
User=dhughes
WorkingDirectory=/home/dhughes/apps/color-the-map
ExecStart=/home/dhughes/apps/color-the-map/venv/bin/python app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
